# systemd/acp-monitor.service
[Unit]
Description=ACP Session Context Monitor
After=network.target

[Service]
Type=oneshot
User=openclaw-user
ExecStart=/opt/health-monitor/scripts/acp-monitor.sh --alert --quiet
EnvironmentFile=/opt/health-monitor/config.env
WorkingDirectory=/opt/health-monitor

[Install]
WantedBy=multi-user.target
